# linux-discord-rpc.service
[Unit]
Description=Discord Rich Presence Daemon

[Service]
ExecStart=rpc-daemon

[Install]
WantedBy=default.target

